# aurbs.timer
[Unit]
Description=Run AurBS build hourly

[Timer]
OnCalendar=*:0,15,30,45
Persistent=true

[Install]
WantedBy=timers.target
