# aurto.service
[Unit]
Description=aurto

[Service]
ExecStart=/opt/aurto.sh
User=unixfox
Group=unixfox

[Install]
WantedBy=multi-user.target
