# config/aurbs/aurbs.timer
[Unit]
Description=Run AurBS build hourly

[Timer]
OnCalendar=*:0/15
Unit=aurbs.service

[Install]
WantedBy=timers.target
